# key.service
[Unit]
Description = key

[Service]
ExecStart = /home/pi/Key/start.sh
Type = simple

[Install]
WantedBy = multi-user.target
